# src-sh/iso-overlay/etc/sv/installer/run
#!/bin/bash
[ -r conf ] && . ./conf
read -r tty < /sys/class/tty/console/active
tty=/dev/${tty##* }
exec setsid sulogin -c start-trident-installer
